# modules/cabine/templates/dev.couchpotato.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>org.atonie.couchpotato</string>

    <key>Boxen</key>
    <dict>
      <key>Description</key>
      <string>Automatic Movie Downloading via NZBs</string>
    </dict>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/python</string>
      <string><%= scope.lookupvar "cabine::couchpotato::executable" %></string>
      <string>--config_file</string>
      <string><%= scope.lookupvar "cabine::couchpotato::configfile" %></string>
      <string>--data_dir</string>
      <string><%= scope.lookupvar "cabine::config::couchpotato_datadir" %></string>
      <string>--console_log</string>
    </array>

    <key>RunAtLoad</key>
    <true />

    <key>KeepAlive</key>
    <false />

    <key>WorkingDirectory</key>
    <string><%= scope.lookupvar "cabine::config::couchpotato_repodir" %></string>

    <key>StandardErrorPath</key>
    <string><%= scope.lookupvar "cabine::config::couchpotato_logdir" %>/console.log</string>

    <key>StandardOutPath</key>
    <string><%= scope.lookupvar "cabine::config::couchpotato_logdir" %>/console.log</string>
</dict>
</plist>
