# modules/people/templates/backup.github.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
"http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>backup.github</string>

  <key>ProgramArguments</key>
  <array>
    <string><%= scope.lookupvar "people::sr::backup_github_bin" %></string>
    <string><%= scope.lookupvar "people::sr::backup_github_dir" %></string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string><%= scope.lookupvar "people::sr::home" %>/bin:/usr/bin:/bin</string>
  </dict>

  <key>StandardOutPath</key>
  <string><%= scope.lookupvar "people::sr::backup_github_log" %></string>
  <key>StandardErrorPath</key>
  <string><%= scope.lookupvar "people::sr::backup_github_log" %></string>

  <key>StartInterval</key>
  <integer><%= scope.lookupvar "people::sr::backup_github_interval" %></integer>
</dict>
</plist>
